# docker/supervisor/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:queue-high]
command=php /var/www/html/artisan queue:work --queue=high --tries=5 --backoff=10 --sleep=3
numprocs=3
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-high.log

[program:queue-default]
command=php /var/www/html/artisan queue:work --queue=default --tries=3 --backoff=30 --sleep=3
numprocs=2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-default.log

[program:queue-low]
command=php /var/www/html/artisan queue:work --queue=low --tries=3 --backoff=60 --sleep=3
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-low.log

[program:scheduler]
command=bash -c "while true; do php /var/www/html/artisan schedule:run; sleep 60; done"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
